# hid.service
[Unit]
Description=Enable HID
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/hid

[Install]
WantedBy=multi-user.target
